# system/root/etc/services.d/git-daemon/run
#!/usr/bin/with-contenv bash
set -Eeuo pipefail

exec as-user bash /opt/git-server/git-daemon.sh
