# system/root/etc/services.d/git-http/run
#!/usr/bin/with-contenv bash
set -Eeuo pipefail

cd /opt/git-server

exec as-user bash /opt/git-server/git-http.sh
