# system/root/etc/services.d/sshd/run
#!/bin/bash
# services.d/sshd/run
exec /usr/sbin/sshd -D -e
